# infra/prod/systemd/momentstudio-backup.service.example
[Unit]
Description=momentstudio backup (DB + media)
Wants=network-online.target
After=network-online.target docker.service
Requires=docker.service

[Service]
Type=oneshot
# Update this path to match your deploy directory (example: /opt/momentstudio).
WorkingDirectory=/opt/momentstudio
ExecStart=/opt/momentstudio/infra/prod/backup.sh

